# root/etc/services.d/local-decrypt/run
#!/usr/bin/with-contenv sh
. "/usr/local/bin/logger"
# ==============================================================================
# local-decrypt
# Mounts the local-decrypt directory (decrypted local files) using rclone
# ==============================================================================

echo "Mounting local-decrypt..." | info

set -x
exec s6-notifyoncheck -n 1000 \
      fdmove 2 1 \
        rclone mount \
          --config=/tmp/rcloneconfig/rclone.conf \
          -v \
          "${RCLONE_LOCAL_DECRYPT_REMOTE}:${RCLONE_LOCAL_DECRYPT_DIR}" \
          /mounts/local-decrypt
